# supervisor/jindai.conf
[program:jindai-api]
directory=/app/
command=gunicorn -w 1 --threads 4 --chdir /app -b 0.0.0.0:8370 --reload api:app
autostart=true
